# container-assets/systemd/discobot-agent-api.service
[Unit]
Description=Discobot agent API (AI coding agent interface)
After=discobot-setup.service discobot-proxy.service
Requires=discobot-setup.service
Wants=discobot-proxy.service discobot-dockerd.service

[Service]
Type=simple
User=discobot
Group=discobot
WorkingDirectory=/home/discobot/workspace
EnvironmentFile=/run/discobot/agent-env
ExecStart=/opt/discobot/bin/discobot-agent-api
Restart=on-failure
RestartSec=5

StandardOutput=journal+console
StandardError=journal+console

[Install]
WantedBy=multi-user.target
