# container-assets/systemd/discobot-dockerd.service
[Unit]
Description=Discobot Docker daemon (nested Docker-in-Docker)
After=discobot-setup.service discobot-proxy.service
Requires=discobot-setup.service
Wants=discobot-proxy.service

[Service]
Type=simple
EnvironmentFile=/run/discobot/container-env
EnvironmentFile=-/run/discobot/proxy-env
ExecStart=/usr/sbin/dockerd --data-root /.data/docker --storage-driver overlay2 --host unix:///var/run/docker.sock --log-level error
ExecStartPost=/bin/bash -c 'for i in $(seq 1 60); do [ -S /var/run/docker.sock ] && chmod 0666 /var/run/docker.sock && exit 0; sleep 0.5; done; exit 1'
Restart=on-failure
RestartSec=5

StandardOutput=journal+console
StandardError=journal+console

[Install]
WantedBy=multi-user.target
