# container-assets/systemd/websockify.service
[Unit]
Description=WebSocket to VNC proxy (websockify)
# Bridges WebSocket connections on port 6080 to VNC on port 5900.
# x11vnc.socket handles on-demand start of x11vnc → xvfb when a connection arrives.

[Service]
Type=simple
ExecStart=/usr/bin/websockify 6080 localhost:5900
Restart=on-failure
RestartSec=2

[Install]
WantedBy=multi-user.target
