# container-assets/systemd/websockify.socket
[Unit]
Description=WebSocket to VNC proxy socket activation on port 6080

[Socket]
ListenStream=6080

[Install]
WantedBy=sockets.target
