# container-assets/systemd/x11vnc.socket
[Unit]
Description=x11vnc VNC server socket activation on port 5900

[Socket]
ListenStream=5900
Accept=true

[Install]
WantedBy=sockets.target
